# test/config/supervisor/conf.d/apache.conf
[program:apache]
autostart=false
command=/usr/sbin/apache2 -d @ROOTDIR@/output/config/apache2 -f @ROOTDIR@/output/config/apache2/apache2.conf -DFOREGROUND

