# test/config/supervisor/conf.d/lighttpd.conf
[program:lighttpd]
autostart=false
command=/usr/sbin/lighttpd -D -f @ROOTDIR@/output/config/lighttpd/lighttpd.conf 

