# test/config/supervisor/conf.d/nginx.conf
[program:nginx]
autostart=false
command=/usr/sbin/nginx -p @ROOTDIR@ -c @ROOTDIR@/output/config/nginx/nginx.conf -g 'error_log stderr;' 

